# cache.service
[Unit]
Description=ProtonVPN Cache

[Service]
ExecStart=~/.cache/cache